# config/lxd-agent/lxd-agent.service
# Generated by distrobuilder
#
[Unit]
Description=LXD Guest Agent
Documentation=https://documentation.ubuntu.com/lxd/en/latest/
Before=multi-user.target cloud-init.target cloud-init.service cloud-init-local.service
DefaultDependencies=no

[Service]
Type=notify
WorkingDirectory=-/run/lxd-agent
ExecStartPre=/usr/lib/systemd/lxd-agent-setup
ExecStart=/run/lxd-agent/lxd-agent
Restart=on-failure
RestartSec=5s
StartLimitInterval=60
StartLimitBurst=10